# src/com.gurjar.MyDownloadsFolderMonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>MyDownloadsFolderMonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>java</string>
        <string>-jar</string>
        <string>/path/to/MyDownloadsFolderMonitor.jar</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/MyDownloadsFolderMonitor.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/MyDownloadsFolderMonitor.err.log</string>
</dict>
</plist>
